# rootfs/etc/profile.d/chamber.sh
export CHAMBER_KMS_KEY_ALIAS="alias/${NAMESPACE}-${ENVIRONMENT}-chamber"
